# bluepill/bluepill.init.d.erb
#! /bin/sh
### BEGIN INIT INFO
# Provides:          bluepill
# Required-Start:
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: bluepill daemon to monitor important processes
# Description:       https://github.com/arya/bluepill
### END INIT INFO

# Author: Tom Sommer

export PATH="/usr/local/bin:$PATH"
 
# Check for missing binaries
BABUSHKA_BIN=/usr/local/bin/babushka
test -x $BABUSHKA_BIN || { echo "$BABUSHKA_BIN not installed";
        if [ "$1" = "stop" ]; then exit 0;
        else exit 5; fi; }

RVM_SETUP='source "$HOME/.rvm/scripts/rvm"'
APP_USER='su - <%= username %> -c'

# Create run dir when necessary
RUNDIR=/var/run/bluepill
if [ ! -d $RUNDIR ]; then
  mkdir -p $RUNDIR
  chmod 777 $RUNDIR
fi

case "$1" in
  start)
    echo -n "Starting bluepill "
    $APP_USER "$RVM_SETUP && babushka tsommer:start-bluepill.task --update"
    ;;
  stop)
    echo -n "Shutting down bluepill "
    $APP_USER "$RVM_SETUP && babushka tsommer:stop-bluepill.task --update"
    ;;
  restart)
    ## Stop the service and regardless of whether it was
    ## running or not, start it again.
    $0 stop
    $0 start
    ;;
  *)
    echo "Usage: $0 {start|stop|restart}"
    exit 1
    ;;
esac

